# etc/init.d/zynthian-webconf
#!/sbin/openrc-run

export PATH=/usr/local/bin:/zynthian/zynthian-sys/sbin:/bin/:/usr/bin:/sbin:/usr/sbin
export HOME=/root
export ZYNTHIAN_DIR=/zynthian
export ZYNTHIAN_SYS_DIR=/zynthian/zynthian-sys
export ZYNTHIAN_CONFIG_DIR=/zynthian/zynthian-webconf
export PYTHONPATH=/usr/lib/python3.11/dist-packages:/zynthian/zynthian-ui

name=$RC_SVCNAME
description="Zynthian Webconf"
supervisor="supervise-daemon"
command="/zynthian/zynthian-webconf/zynthian_webconf.sh"
output_log="/var/log/zynthian_webconf.log"
error_log="/var/log/zynthian_webconf.err"