# etc/init.d/zynthian-webconf
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Zynthian Webconf"
supervisor="supervise-daemon"
command="/bin/bash"
command_args="#ZYNTHIAN_DIR#/zynthian-webconf/zynthian_webconf.sh"
output_log="/var/log/zynthian_webconf.log"
error_log="/var/log/zynthian_webconf.err"
